# host/unsloth/unsloth-studio.service.d/override.conf
# Drop-in override (OURS) — pins the served model, quant, context window,
# concurrency, and KV-cache precision. We run the text-only unsloth/Qwen3-8B-GGUF
# (UD-Q4_K_XL) with an 80k total context split across 2 parallel slots (40k tokens
# per slot — the model's native 40,960 window) and a q8_0-quantized KV cache.
#
# Why these values:
# - --max-seq-length 81920 sets llama.cpp's TOTAL -c; with --parallel 2 each
#   concurrent slot gets 81920/2 = 40960 tokens — EXACTLY the model's native
#   context (the GGUF's qwen3.context_length=40960), so NO rope/YaRN scaling is
#   needed and there's no long-context quality degradation. 40k/slot comfortably
#   holds OpenWebUI's web-search results (full page text is injected, not
#   summarized), which overflowed the first 16k/slot (-c 32768) window and were
#   tight at 32k/slot. Going beyond 40k/slot would require YaRN (or the separate
#   unsloth/Qwen3-8B-128K-GGUF variant) AND more VRAM — not worth it here.
# - --cache-type-k/-v q8_0 stores the KV cache at 8-bit instead of f16, roughly
#   HALVING its VRAM. That keeps the whole LLM at ~11 GB (weights ~4.8 GB + ~6 GB
#   q8 KV), leaving ~12 GB of the RTX 3090's 24 GB free for ComfyUI and other GPU
#   workloads. q8_0 KV is effectively lossless for chat/tool-use, and requires
#   flash-attn (Studio enables it by default).
# - No --mmproj: Qwen3-8B is text-only, so no vision projector is loaded.
# - --chat-template-kwargs '{"enable_thinking": false}': Studio defaults thinking
#   ON for Qwen3-8B, but the <think> blocks slow replies and break n8n AI-Agent
#   tool-call parsing, so we disable it (last-wins over Studio's auto-added
#   enable_thinking:true). The JSON is single-quoted so systemd keeps it as one
#   argv token and preserves its inner double-quotes.
#
# Studio only MANAGES the model/networking/auth flags; every other llama-server
# flag (-c, --parallel, --cache-type-*, --chat-template-kwargs, ...) passes
# straight through, last-wins (see Studio's llama_server_args.py denylist).
# The empty ExecStart= clears the base ExecStart before re-setting it (systemd
# requires this for Type=simple single-ExecStart units).
# On jarvis: /etc/systemd/system/unsloth-studio.service.d/override.conf
[Service]
ExecStart=
ExecStart=/opt/unsloth/.unsloth/studio/unsloth_studio/bin/unsloth run -H 0.0.0.0 -p 8888 --model unsloth/Qwen3-8B-GGUF --gguf-variant UD-Q4_K_XL --max-seq-length 81920 --parallel 2 --chat-template-kwargs '{"enable_thinking": false}' --cache-type-k q8_0 --cache-type-v q8_0
