# host/unsloth/unsloth-studio.service.d/override.conf
# Drop-in override (OURS) — raise the served context window to the model's
# native 256k. Unsloth Studio's --max-seq-length defaults to 65536; the model +
# llama.cpp support 262144, and llama.cpp already allocates the 256k KV cache, so
# this costs no extra VRAM. The empty ExecStart= clears the base ExecStart before
# re-setting it (systemd requires this for Type=simple single-ExecStart units).
# On jarvis: /etc/systemd/system/unsloth-studio.service.d/override.conf
[Service]
ExecStart=
ExecStart=/opt/unsloth/.unsloth/studio/unsloth_studio/bin/unsloth run -H 0.0.0.0 -p 8888 --model unsloth/Qwen3.5-4B-GGUF --max-seq-length 262144
