# host/unsloth/unsloth-studio.service.d/override.conf
# Drop-in override (OURS) — pins the served model, quant, context window, and
# concurrency. We run the text-only unsloth/Qwen3-8B-GGUF (UD-Q4_K_XL) at a 32k
# context with 2 parallel slots. The small context deliberately keeps the KV
# cache tiny so the whole LLM fits in ~8 GB, leaving the bulk of the RTX 3090's
# 24 GB free for ComfyUI and other GPU workloads. --max-seq-length sets
# llama.cpp's total -c; --parallel passes straight through to llama-server
# (Studio only manages the model/networking/auth flags — every other flag is
# last-wins pass-through). No --mmproj: Qwen3-8B is text-only, so no vision
# projector is loaded. We also force --chat-template-kwargs '{"enable_thinking":
# false}': Studio defaults thinking ON for Qwen3-8B, but the <think> blocks slow
# replies and break n8n AI-Agent tool-call parsing, so we disable it (last-wins
# over Studio's auto-added enable_thinking:true). The JSON is single-quoted so
# systemd keeps it as one argv token and preserves its inner double-quotes.
# The empty ExecStart= clears the base ExecStart before re-setting it (systemd
# requires this for Type=simple single-ExecStart units).
# On jarvis: /etc/systemd/system/unsloth-studio.service.d/override.conf
[Service]
ExecStart=
ExecStart=/opt/unsloth/.unsloth/studio/unsloth_studio/bin/unsloth run -H 0.0.0.0 -p 8888 --model unsloth/Qwen3-8B-GGUF --gguf-variant UD-Q4_K_XL --max-seq-length 32768 --parallel 2 --chat-template-kwargs '{"enable_thinking": false}'
